# .docker/config/supervisor.conf
[supervisord]
nodaemon=true

[program:sysctl]
command=sysctl -p